# docker/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
childlogdir=/var/log/supervisor
pidfile=/var/run/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:langflow]
command=langflow run --host %(ENV_LANGFLOW_HOST)s --port %(ENV_LANGFLOW_PORT)s
priority=10
autostart=true
autorestart=true
startretries=3
stopsignal=TERM
stdout_logfile=/var/log/langflow.log
stderr_logfile=/var/log/langflow_error.log

[program:nginx]
command=nginx -g "daemon off;"
priority=20
autostart=true
autorestart=true
startretries=3
stopsignal=TERM
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx_error.log